# scripts/start_webservice.conf
description "Start webservice"

start on runlevel [2345]

exec /root/webservice